# alpine_cnode_scripts_and_services/etc/init.d/grafana
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Grafana service"

depend() {
        after network-online
}

start() {

        ebegin "Starting $RC_SVCNAME"
        start-stop-daemon --background --start --exec /usr/sbin/grafana-server \
        -- --homepath=/usr/share/grafana \
        --config=/etc/grafana.ini \
        --make-pidfile --pidfile /var/run/grafana.pid \
        eend $?

}

start_post() {
rcstatus=blank
        while [ -f /var/run/grafana.pid ]
        do
                sleep 10
                rcstatus=`rc-service $RC_SVCNAME status | awk '{print $NF}'`
                if [ -z $rcstatus ]; then
                        rc-service -c $RC_SVCNAME restart
                        break
                fi
        done &
}

stop() {
        ebegin "Stopping $RC_SVCNAME"
        start-stop-daemon --stop --exec /usr/sbin/grafana-server \
        --pidfile /var/run/grafana.pid \
        -s 2
        eend $?
}
